# docker-nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:scampi]
command=/usr/scampi.sh

[program:scampiREST]
command=/usr/scampiRest.sh